# bitcoin/contrib/init/org.h2hcoin.h2hcoind.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.h2hcoin.h2hcoind</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/h2hcoind</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
